# files/systemd/ansible-pull.timer
[Unit]
Description='Runs ansible-pull daily and on vm boot'

[Timer]
OnBootSec=2min
OnUnitActiveSec=1d 

[Install]
WantedBy=timers.target

